# package/base-files/files/etc/uci-defaults/14_prplos_defaults
#!/bin/sh

. /lib/functions/uci-defaults.sh

force_nec_wx3000hp_wireless_defaults() {
	local ugw852_wireless_md5sum="a243e630fc195e7ec45010f80d821317"
	[ ! -f /etc/config/wireless ] && return
	[ "$ugw852_wireless_md5sum" != "$(md5sum /etc/config/wireless | cut -d' ' -f1)" ] && return

	logger -t "uci-defaults" -p daemon.warn "Forcing prplOS wireless defaults"
	cat > "/etc/config/wireless" <<EOF
config wifi-device 'radio0'
        option phy 'phy0'
        option macaddr '02:9A:96:FB:59:0F'
        option type 'mac80211'
        option band '2.4GHz'
        option channel 'auto'
        option hwmode '11bgnax'
        option htmode 'HT20'
        option disabled '1'
        option country 'US'
        option num_antennas '4'
        option beacon_int '100'
        option short_gi_20 '1'
        option short_gi_40 '1'
        option txpower '100'
        option sDisableMasterVap '1'
        option multi_ap_profile '2'

config wifi-iface 'default_radio100'
        option device 'radio0'
        option ifname 'wlan0'
        option macaddr '02:9A:96:FB:59:0F'
        option network 'lan'
        option mode 'ap'
        option hidden '1'
        option ssid 'dummy_ssid_0'
        option start_disabled '1'

config wifi-iface 'default_radio10'
        option device 'radio0'
        option ifname 'wlan0.0'
        option macaddr '02:9A:96:FB:59:11'
        option network 'lan'
        option mode 'ap'
        option encryption 'none'
        option wav_security_mode 'None'
        option isolate '0'
        option uapsd '1'
        option hidden '0'
        option wmm '1'
        option macfilter 'disable'
        option ieee80211w '1'
        option max_inactivity '60'
        option wav_beacon_type 'None'
        option wps_pushbutton '1'
        option wps_rf_bands 'ag'
        option ssid 'prplOS'
        option rrm_neighbor_report '1'
        option start_disabled '1'

config wifi-iface 'default_radio11'
        option device 'radio0'
        option ifname 'wlan0.1'
        option macaddr '02:9A:96:FB:59:12'
        option network 'guest'
        option mode 'ap'
        option encryption 'none'
        option wav_security_mode 'None'
        option isolate '0'
        option uapsd '1'
        option hidden '0'
        option wmm '1'
        option macfilter 'disable'
        option ieee80211w '1'
        option max_inactivity '60'
        option wav_beacon_type 'None'
        option wps_pushbutton '1'
        option wps_rf_bands 'ag'
        option ssid 'prplOS-guest'
        option rrm_neighbor_report '1'
        option start_disabled '1'

config wifi-device 'radio2'
        option phy 'phy1'
        option macaddr '02:9A:96:FB:59:20'
        option type 'mac80211'
        option band '5GHz'
        option channel 'auto'
        option hwmode '11anacax'
        option htmode 'VHT80'
        option country 'US'
        option disabled '1'
        option num_antennas '4'
        option beacon_int '100'
        option doth '1'
        option short_gi_20 '1'
        option short_gi_40 '1'
        option txpower '100'
        option sDisableMasterVap '1'
        option multi_ap_profile '2'

config wifi-iface 'default_radio102'
        option device 'radio2'
        option ifname 'wlan2'
        option macaddr '02:9A:96:FB:59:20'
        option network 'lan'
        option mode 'ap'
        option hidden '1'
        option ssid 'dummy_ssid_2'
        option start_disabled '1'

config wifi-iface 'default_radio42'
        option device 'radio2'
        option ifname 'wlan2.0'
        option macaddr '02:9A:96:FB:59:22'
        option network 'lan'
        option mode 'ap'
        option encryption 'none'
        option wav_security_mode 'None'
        option isolate '0'
        option uapsd '1'
        option hidden '0'
        option wmm '1'
        option macfilter 'disable'
        option ieee80211w '1'
        option max_inactivity '60'
        option wav_beacon_type 'None'
        option wps_pushbutton '1'
        option wps_rf_bands 'ag'
        option ssid 'prplOS'
        option rrm_neighbor_report '1'
        option start_disabled '1'

config wifi-iface 'default_radio43'
        option device 'radio2'
        option ifname 'wlan2.1'
        option macaddr '02:9A:96:FB:59:23'
        option network 'guest'
        option mode 'ap'
        option encryption 'none'
        option wav_security_mode 'None'
        option isolate '0'
        option uapsd '1'
        option hidden '0'
        option wmm '1'
        option macfilter 'disable'
        option ieee80211w '1'
        option max_inactivity '60'
        option wav_beacon_type 'None'
        option wps_pushbutton '1'
        option wps_rf_bands 'ag'
        option ssid 'prplOS-guest'
        option rrm_neighbor_report '1'
        option start_disabled '1'
EOF
}

board=$(board_name)

case "$board" in
"EASY350 ANYWAN (GRX350) Main model" |\
"EASY350 ANYWAN (GRX350) Axepoint Asurada model")
	force_nec_wx3000hp_wireless_defaults
	;;
esac

exit 0
